# daemon/auto-patcher.timer
[Unit]
Description=Trigger Auto-Patcher Service

[Timer]
Unit=auto-patcher.service
OnCalendar=daily

[Install]
WantedBy=multi-user.target